# deploy/systemd/spotifyd.root.service
[Unit]
Description=spotifyd headless Spotify Connect daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=HOME=/root
ExecStartPre=/bin/sh -lc 'pulseaudio --check >/dev/null 2>&1 || pulseaudio --start'
ExecStart=/usr/bin/env spotifyd --no-daemon --config-path /root/.config/spotifyd/spotifyd.conf
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
