# deploy/systemd/spotifyd.user.service
[Unit]
Description=spotifyd headless Spotify Connect daemon
After=network-online.target pulseaudio.service
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/sh -lc 'pulseaudio --check >/dev/null 2>&1 || pulseaudio --start'
ExecStart=/usr/bin/env spotifyd --no-daemon --config-path %h/.config/spotifyd/spotifyd.conf
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
