# deploy/systemd/telegram-spotify-bot.root.service
[Unit]
Description=Telegram Spotify controller bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/spotify247
Environment=HOME=/root
EnvironmentFile=/root/spotify247/.env
ExecStart=/root/spotify247/.venv/bin/python -m app.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
